# etc/profile.d/labl.sh
PATH=${PATH}:/opt/labl/bin:/opt/labl/sbin
